# deploy/sms-scheduler.timer
[Unit]
Description=SMS Admin Scheduler Timer
Documentation=man:systemd.timer(5)

[Timer]
# Start 30 seconds after boot
OnBootSec=30
# Then run every 60 seconds
OnUnitActiveSec=60
# Randomize delay up to 5 seconds to prevent thundering herd
RandomizedDelaySec=5
# Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target
