# deploy/sms-scheduler.service
[Unit]
Description=SMS Admin Scheduler Service
After=network.target

[Service]
User=smsadmin
Group=smsadmin
WorkingDirectory=/opt/sms-admin
Environment="PATH=/opt/sms-admin/venv/bin"
EnvironmentFile=/opt/sms-admin/.env
ExecStartPre=/usr/local/bin/dbdoctor --apply
ExecStart=/opt/sms-admin/venv/bin/python -m app.scheduler_runner
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
