# deploy/sms-scheduler.timer
[Unit]
Description=SMS Admin Scheduler Timer
Documentation=man:systemd.timer(5)

[Timer]
# Start 30 seconds after boot
OnBootSec=30
# Run every 30 seconds for timely message delivery
OnUnitActiveSec=30
# Persist timer state across reboots (catches up missed runs)
Persistent=true

[Install]
WantedBy=timers.target
